# backend/piafadex.service
[Unit]
Description=PIAFADEX FastAPI backend
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/piafadex/backend
EnvironmentFile=/opt/piafadex/backend/.env
ExecStart=/opt/piafadex/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
